# projectTrend/.crontab
PATH=/bin:/usr/bin
0 7,13 * * * cd /local/scratch/philippeg/jira/projectTrend;./genreport.sh;./deploy.sh


